# pine64/pinephone/etc/init/set-info.conf
description "Update machine info"

task

start on startup

script
    hostnamectl --pretty set-chassis handset
end script
